# deploy/systemd/smart-display-kiosk.service
[Unit]
Description=Smart Display kiosk session
After=network-online.target smart-display.service
Requires=smart-display.service

[Service]
Type=simple
User=pi
Group=pi
ExecStart=/usr/bin/xinit /opt/smart-display/deploy/x11/kiosk-session.sh -- :0 vt1 -nolisten tcp
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

