# deploy/systemd/smart-display.service
[Unit]
Description=Smart Display backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/smart-display
EnvironmentFile=/opt/smart-display/.env
ExecStart=/opt/smart-display/.venv/bin/python -m smart_display.app
Restart=always
RestartSec=3
User=pi
Group=pi
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

